# infra/systemd/elamora-backup-failure-alert.service
[Unit]
Description=Send Elamora backup failure alert
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/debian/apps/elamora
EnvironmentFile=/home/debian/apps/elamora/server/.env
ExecStart=/usr/bin/node /home/debian/apps/elamora/server/dist/scripts/send-backup-status-email.js failure
