# infra/systemd/elamora-backup-daily-check.service
[Unit]
Description=Check daily Elamora database backup
After=network-online.target elamora-db-backup.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/debian/apps/elamora
EnvironmentFile=/home/debian/apps/elamora/server/.env
Environment=BACKUP_REPORT_EMAIL=resulramadan@icloud.com
ExecStart=/usr/bin/node /home/debian/apps/elamora/server/dist/scripts/send-backup-status-email.js daily
SuccessExitStatus=2
