# infra/systemd/elamora-pending-payment-reminders.timer
[Unit]
Description=Run Elamora pending payment admin reminders daily

[Timer]
OnCalendar=*-*-* 09:00:00 Europe/Belgrade
Persistent=true
Unit=elamora-pending-payment-reminders.service

[Install]
WantedBy=timers.target
