# infra/systemd/elamora-db-backup.service.d/alert.conf
[Unit]
OnFailure=elamora-backup-failure-alert.service
