# infra/systemd/elamora-db-backup.service.d/environment.conf
[Service]
EnvironmentFile=/home/debian/apps/elamora/server/.env
